# supervisor/listener.ini
[program:listener]
directory=/home/ec2-user/listener-app/src
command=/home/ec2-user/listener-app/src/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/listener.err.log
stdout_logfile=/var/log/listener.out.log
environment=CONFIG_FILE="/home/ec2-user/config.yml"